# deploy/systemd/cms-data-status.timer
[Unit]
Description=Daily CMS public-data publisher freshness discovery

[Timer]
OnCalendar=*-*-* 06:15:00 UTC
RandomizedDelaySec=15m
Persistent=true
AccuracySec=1m
Unit=cms-data-status.service

[Install]
WantedBy=timers.target
